# PF/PB/TFC/crontab
53 23 * * * * chmod +x ./makepayments.sh
54 23 * * * * ./makepayments.sh
